# var/lib/dpkg/info/dbus.list
/etc
/etc/init.d
/etc/init.d/dbus
/etc/init
/etc/init/dbus.conf
/etc/default
/etc/default/dbus
/var
/var/lib
/var/lib/dbus
